# files/etc/profile.d/ontoportal.sh
#useful aliases for OntoPortal
alias suo="sudo su - ontoportal"
alias gtc="pushd /opt/ontoportal/ncbo_cron"
alias gta="pushd /opt/ontoportal/ontologies_api/current"
alias gtu="pushd /opt/ontoportal/bioportal_web_ui/current"
alias gtv="pushd /opt/ontoportal/virtual_appliance/"
alias tla="tail -f /opt/ontoportal/ontologies_api/current/log/appliance.log"
alias tlu="tail -f /opt/ontoportal/bioportal_web_ui/current/log/appliance.log"
alias tlc="tail -f /opt/ontoportal/ncbo_cron/logs/scheduler.log"
